# build/rootfs/etc/systemd/system/oxware-installer.service
[Unit]
Description=OXware Live Installer
After=network.target getty@tty1.service
Before=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/oxware-installer/install.py
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
